# newapi-monitor-once@.service
[Unit]
Description=New-API one-shot token group poll (%i)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/root/projects/zhongzhuan
ExecStart=/root/projects/zhongzhuan/.venv/bin/python \
  /root/projects/zhongzhuan/newapi_monitor.py \
  --env-file /root/projects/zhongzhuan/sites/%i.env
TimeoutStartSec=240
UMask=0077
Nice=10
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/root/projects/zhongzhuan/data
